# systemdExamples/iab.service
[Unit]
Description=IncusAutobackup (IAB)

[Service]
Type=oneshot
WorkingDirectory=/root/
ExecStart=/root/iab